# src/daemon/rpm-ostreed.service.in
[Unit]
Description=OSTree+RPM package management daemon
ConditionPathExists=/ostree

[Service]
Type=dbus
BusName=org.projectatomic.rpmostree1
NotifyAccess=main
@SYSTEMD_ENVIRON@
ExecStart=@bindir@/rpm-ostree start-daemon
ExecReload=@libexecdir@/@primaryname@ reload
