# src/daemon/rpm-ostreed.service.in
[Unit]
Description=RPM OSTree Manager

[Service]
Type=dbus
BusName=org.projectatomic.rpmostree1
ExecStart=@libexecdir@/rpm-ostreed
